# etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

/usr/sbin/nginx 2>&1 | awk '{print "[webapp] " $0}'

